# ansible/roles/vaultCluster/templates/vaultFollowerConf.service.j2
[Service]
Environment="VAULT_TOKEN={{ root_key.stdout }}"